# desktop/claude-code-hygiene/stale-claude-reap.user.service
[Unit]
Description=Stale (terk edilmis, >24h idle) Claude Code session reaper (swap geri kazan, idle-safe, graceful)

[Service]
Type=oneshot
ExecStart=%h/.claude/scripts/stale-claude-reap.sh 24
Nice=10
IOSchedulingClass=idle
