# desktop/claude-code-hygiene/stale-claude-reap.user.timer
[Unit]
Description=stale-claude-reap timer (her 6 saatte bir >24h idle session kontrolu)

[Timer]
OnBootSec=20min
OnUnitActiveSec=6h
Persistent=true
RandomizedDelaySec=3min

[Install]
WantedBy=timers.target
